# usb-gadget.service
[Unit]
Description=RPi0 USB Gadget

[Service]
Type=forking
RemainAfterExit=yes
ExecStart=/sbin/usb-gadget-init

[Install]
WantedBy=sysinit.target
